# cmd/install/bin/sg-init.service
[Unit]
Description=Performs the initial sourcegraph setup setups.

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/local/bin/sg-init

[Install]
WantedBy=multi-user.target
