# cmd/install/bin/sourcegraphd.service
[Unit]
Description=Monitors sourcegraph install and performs updates.
After=network.target
After=sg-init.service

[Service]
ExecStart=/usr/local/bin/sourcegraphd
RestartSec=3
KillMode=process

[Install]
WantedBy=multi-user.target
